# modules/addon-sendchanges/templates/crontab
MAILTO=release@mozilla.com

0 6 * * 6   <%= addonSendchangeUser %> source <%= addonProfile %> && bash <%= sendchangeScript %> -p <%= addonSendchangePort %> -r -a <%= addonsFile %> -m <%= fqdn %>
