# modules/puppet/templates/com.reductivelabs.puppet.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>EnvironmentVariables</key>
        <dict>
                <key>PATH</key>
                <string>/sbin:/usr/sbin:/bin:/usr/bin</string>
                <key>RUBYLIB</key>
                <string>/usr/lib/ruby/site_ruby/1.8/:/usr/lib/ruby/user-gems/1.8/gems</string>
        </dict>
        <key>Label</key>
        <string>com.reductivelabs.puppet</string>
        <key>OnDemand</key>
        <false/>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/local/bin/sleep-and-run-puppet.sh</string>
                <string><%=puppetServer%></string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>ServiceDescription</key>
        <string>Puppet Daemon</string>
        <key>ServiceIPC</key>
        <false/>
        <key>StandardErrorPath</key>
        <string>/var/puppet/log/puppet.err</string>
        <key>StandardOutPath</key>
        <string>/var/puppet/log/puppet.out</string>
</dict>
</plist>
