# templates/stage.root.crontab.erb
MAILTO=release@mozilla.com

1 * * * *	find /data/ftp/pub/{firefox,xulrunner,mobile}/tinderbox-builds/ -mindepth 2 -maxdepth 2 -type d -mmin +43200 -name 1????????? -exec rm -rf {} \;
6 * * * *	find /data/ftp/pub/firefox/try-builds/ -mindepth 1 -maxdepth 1 -type d -mmin +43200 -exec rm -rf {} \;
<% if scope.lookupvar('deleteCandidatesDirs') == '1' %>
11 * * * *	find /data/ftp/pub/{firefox,xulrunner,mobile}/nightly/ -mindepth 2 -maxdepth 2 -name build? -type d -mtime +1 -exec rm -rf {} \;
<% end %>
16 * * * *	find /data/ftp/pub/{firefox,xulrunner,mobile}/nightly/20??/ -mindepth 2 -maxdepth 2 -mtime +1 -exec rm -rf {} \;
21 * * * *	find /data/ftp/pub/{firefox,xulrunner,mobile}/nightly/ -type l -mtime +1 -exec rm -rf {} \;
26 * * * *	find /data/ftp/pub/{firefox,xulrunner,mobile}/nightly/latest-* -type f -mtime +1 -exec rm -f {} \;
31 * * * *	find /data/symbols_*/ -mindepth 2 -maxdepth 2 -type d -mmin +21600 -exec rm -rf {} \;
36 * * * *	find /data/symbols_*/ -maxdepth 0 -mmin +21600 -type f -exec rm -rf {} \;
